# site/profiles/files/Fedora/config/19/puppet.cron
10,40 * * * * root sleep $((RANDOM % 1500)) && puppet agent -t > /dev/null 2>&1

