# site/profiles/files/Fedora/config/19/puppet.cron
05,35 * * * * root /usr/local/bin/stagger_puppet

