# site/profiles/files/Fedora/config/21/forge.cron
01 * * * * root /usr/local/bin/grub2
10 * * * * root /usr/local/bin/nmt-dracut --quiet
40 * * * * root /usr/local/bin/nmt-dracut --quiet

