# site/profiles/files/Fedora/config/21/media-winpart.mount
[Unit]
Description=Windows partition

[Mount]
What=/dev/sda2
Where=/media/winpart
Type=ntfs
Options=rw
TimeoutSec=5
