# site/profiles/files/Fedora/config/21/puppet.cron
10,40 * * * * root sleep $((RANDOM \% 1500)) && puppet agent -t --environment master > /dev/null 2>&1

