# site/profiles/files/Debian/config/8/ldmagain.service
[Service]
Type=simple
ExecStart=/usr/local/bin/ldmagain
Restart=no
